# Arch/systemd/system/pacman-download-packageupdates.timer
[Unit]
Description=Pacman predownload packages

[Timer]
OnBootSec=1min
OnUnitActiveSec=3h
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
