# Applications/autostart/ckb.desktop
[Desktop Entry]
Type=Application
Name=ckb
Exec=ckb --background
Terminal=false

