# Arch/systemd/user/git-fetch.timer
[Unit]
Description=Fetch git folders in ~/git every 100min

[Timer]
OnStartupSec=30sec
OnUnitActiveSec=50min
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
