# Arch/systemd/user/git-fetch.timer
[Unit]
Description=Fetch git repositories in ~/git

[Timer]
OnCalendar=04/12:42:00
Persistent=true
RandomizedDelaySec=8min

[Install]
WantedBy=timers.target
